# deploy/systemd/golf-grading-sweep.service
[Unit]
Description=Golf Model completed-event grading sweep
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/golf-model
Environment=GOLF_APP_ROOT=/opt/golf-model
Environment=GOLF_DATA_DIR=/opt/golf-model/data
EnvironmentFile=-/opt/golf-model/.env
Nice=10
MemoryMax=3G
ExecStart=/opt/golf-model/venv/bin/python /opt/golf-model/scripts/grading_sweep.py --year 2026 --json

[Install]
WantedBy=multi-user.target
