# deploy/systemd/golf-live-refresh.service
[Unit]
Description=Golf Model Live Refresh Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/golf-model
Environment=PATH=/opt/golf-model/venv/bin:/usr/bin:/bin
Environment=GOLF_APP_ROOT=/opt/golf-model
Environment=GOLF_DATA_DIR=/opt/golf-model/data
EnvironmentFile=-/opt/golf-model/.env
ExecStart=/opt/golf-model/venv/bin/python -m workers.live_refresh_worker
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
