# deploy/systemd/golf-retention.timer
[Unit]
Description=Run Golf Model retention cycle weekly (Sun 04:00 UTC)

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
